# auto-mp4.path
[Unit]
Description="Watch video library for changes"

[Path]
PathChanged=$videodir
Unit=auto-mp4.service

[Install]
WantedBy=multi-user.target
