# auto-mp4.service
[Unit]
Description="Run auto-mp4 script"

[Service]
ExecStart=/usr/local/bin/auto-mp4.sh
